# config/systemd/user/wayland/vm/clip-v2h-save.service
[Unit]
Description=VM Wayland clipboard -> file (v2h.txt)
After=graphical-session.target
Wants=graphical-session.target
ConditionPathExistsGlob=%t/wayland-*

[Service]
Environment=XDG_RUNTIME_DIR=%t
ExecStart=%h/.local/bin/vm_wl_clip_to_file.sh
Restart=always

[Install]
WantedBy=default.target

